# deploy/systemd/fusion-memory-embedding@.service
[Unit]
Description=Fusion Memory embedding worker (%i)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/fusion-memory/embedding-%i.env
ExecStart=%h/.local/bin/fusion-memory embedding-server --host ${FUSION_MEMORY_MODEL_HOST} --port ${FUSION_MEMORY_MODEL_PORT} --model ${FUSION_MEMORY_MODEL_PATH} --device ${FUSION_MEMORY_MODEL_DEVICE} --max-concurrency ${FUSION_MEMORY_MODEL_MAX_CONCURRENCY}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
